# dist/opt/flight/etc/profile.d/80-sshkey.csh
################################################################################
##
## Flight Starter
## Copyright (c) 2019-present Alces Flight Ltd
##
################################################################################
if ( -x "$flight_ROOT/libexec/flight-starter/setup-sshkey" ) then
    $flight_ROOT/libexec/flight-starter/setup-sshkey
endif
